# supervisord.conf
[supervisord]
user=root
nodaemon=true
loglevel=debug

[program:nginx]
command=bash /app/nginx-entrypoint.sh
autorestart=true

[program:n8n]
command=bash /app/n8n-entrypoint.sh
autorestart=true