# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

[program:n8n]
command=bash /app/docker-entrypoint.sh
autorestart=true