# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:nginx]
command=bash /app/nginx-entrypoint.sh
autorestart=true

[program:n8n]
command=bash /app/docker-entrypoint.sh
autorestart=true